# market_maker/bin/crontab
* * * * *  sh /home/py3/daemon.sh
* * * * * sleep 10; sh /home/py3/daemon.sh
* * * * * sleep 20; sh /home/py3/daemon.sh
* * * * * sleep 30; sh /home/py3/daemon.sh
* * * * * sleep 40; sh /home/py3/daemon.sh
* * * * * sleep 50; sh /home/py3/daemon.sh